# recipes-ttt-rootfs/init/files/st-tsn.service
[Unit]
Description=TSN service
After=network.target systemd-networkd.service

[Service]
Type=forking
ExecStart=/usr/sbin/ttt-ip-init-systemd.sh start
ExecStop=/usr/sbin/ttt-ip-init-systemd.sh stop

[Install]
WantedBy=multi-user.target
